# splunk/apps/soclab_detections/default/savedsearches.conf
# =====================================================================
# SOC Lab — 10 SPL Correlation Searches
# Mapped to 12 MITRE ATT&CK techniques across Initial Access, Execution,
# Persistence, Priv Esc, Defense Evasion, Credential Access, Discovery,
# Lateral Movement, and Command & Control.
# =====================================================================

# ---------------------------------------------------------------------
# [1] T1059.001 - PowerShell Encoded Command (Execution)
# Detects PowerShell invoked with -EncodedCommand / -enc / -e flags.
# ---------------------------------------------------------------------
[SOCLab - T1059.001 - PowerShell Encoded Command]
description = PowerShell with encoded command flag (T1059.001). Common in delivery/staging of malicious payloads.
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = index=sysmon EventCode=1 (Image="*\\powershell.exe" OR Image="*\\pwsh.exe") \
  CommandLine IN ("*-enc *", "*-Enc *", "*-EncodedCommand *", "*-encodedcommand *", "* -e *") \
  | eval mitre_technique="T1059.001", mitre_tactic="Execution" \
  | table _time, host, User, Image, CommandLine, ParentImage, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [2] T1059.003 - Suspicious cmd.exe Chains (Execution)
# ---------------------------------------------------------------------
[SOCLab - T1059.003 - Suspicious Cmd Chains]
description = cmd.exe invoked with piped/chained obfuscation (T1059.003).
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = index=sysmon EventCode=1 Image="*\\cmd.exe" \
  CommandLine IN ("*whoami*", "*net user*", "*net localgroup*", "*findstr*password*", "*for /L*", "*powershell*iex*", "*certutil*-decode*") \
  NOT ParentImage IN ("*\\OneDriveSetup.exe", "*\\vmware*") \
  | eval mitre_technique="T1059.003", mitre_tactic="Execution" \
  | table _time, host, User, CommandLine, ParentImage, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [3] T1547.001 - Registry Run Key Persistence (Persistence)
# Sysmon Event 13: Registry value set under Run/RunOnce.
# ---------------------------------------------------------------------
[SOCLab - T1547.001 - Registry Run Key Persistence]
description = Modification to a Run/RunOnce autorun key (T1547.001).
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = index=sysmon EventCode=13 \
  TargetObject IN ("*\\CurrentVersion\\Run\\*", "*\\CurrentVersion\\RunOnce\\*", \
                    "*\\Explorer\\Shell Folders\\Startup*") \
  | eval mitre_technique="T1547.001", mitre_tactic="Persistence" \
  | table _time, host, User, Image, TargetObject, Details, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [4] T1053.005 - Scheduled Task Creation (Persistence / Priv Esc)
# Windows Event 4698 — scheduled task created.
# ---------------------------------------------------------------------
[SOCLab - T1053.005 - Scheduled Task Created]
description = New scheduled task created (T1053.005). Common persistence + priv esc mechanism.
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = (index=wineventlog EventCode=4698) \
  OR (index=sysmon EventCode=1 Image="*\\schtasks.exe" CommandLine="*/create*") \
  | eval taskname=coalesce(TaskName, "unknown") \
  | eval mitre_technique="T1053.005", mitre_tactic="Persistence" \
  | table _time, host, User, Image, CommandLine, taskname, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [5] T1003.001 - LSASS Memory Access (Credential Access)
# Sysmon Event 10: process access to lsass.exe.
# ---------------------------------------------------------------------
[SOCLab - T1003.001 - LSASS Memory Access]
description = Suspicious process opened handle to LSASS with high access mask (T1003.001). Mimikatz-style.
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = index=sysmon EventCode=10 TargetImage="*\\lsass.exe" \
  GrantedAccess IN ("0x1010", "0x1410", "0x1438", "0x143a", "0x1fffff") \
  NOT SourceImage IN ("*\\svchost.exe", "*\\wininit.exe", "*\\services.exe", "*\\csrss.exe", "*\\MsMpEng.exe") \
  | eval mitre_technique="T1003.001", mitre_tactic="Credential Access" \
  | table _time, host, SourceImage, SourceProcessId, GrantedAccess, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [6] T1110 - Brute Force Authentication (Credential Access)
# Repeated failed logon events on Windows (4625) or Linux (sshd failed).
# ---------------------------------------------------------------------
[SOCLab - T1110 - Brute Force Authentication]
description = 10+ failed logons within 5 minutes from same source (T1110). Windows 4625 + Linux SSH failures.
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = (index=wineventlog EventCode=4625) \
  OR (index=linux_auditd sourcetype="linux:auth" "Failed password") \
  | eval src_ip=coalesce(src_ip, IpAddress, host) \
  | eval account=coalesce(TargetUserName, user, "unknown") \
  | stats count as failures, values(account) as accounts by src_ip, host \
  | where failures >= 10 \
  | eval mitre_technique="T1110", mitre_tactic="Credential Access" \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [7] T1070.004 - Windows Event Log Cleared (Defense Evasion)
# Windows Event 1102 (Security log cleared) / 104 (System log cleared).
# ---------------------------------------------------------------------
[SOCLab - T1070.001 - Event Log Cleared]
description = Windows Event Log manually cleared (T1070.001). Strong indicator of anti-forensics.
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = index=wineventlog (EventCode=1102 OR EventCode=104) \
  | eval mitre_technique="T1070.001", mitre_tactic="Defense Evasion" \
  | table _time, host, user, EventCode, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [8] T1021.002 - SMB/Admin Share Lateral Movement (Lateral Movement)
# Windows Event 4624 type 3 (network) + special privileges (4672).
# ---------------------------------------------------------------------
[SOCLab - T1021.002 - Admin Share Lateral Movement]
description = Network logon with admin privileges — potential lateral movement via C$/ADMIN$ (T1021.002).
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = index=wineventlog EventCode=4624 Logon_Type=3 \
  NOT Account_Name IN ("SYSTEM", "ANONYMOUS LOGON", "-") \
  | rex "Source Network Address:\s+(?<src_ip>[^\s]+)" \
  | where isnotnull(src_ip) AND src_ip!="-" \
  | eval mitre_technique="T1021.002", mitre_tactic="Lateral Movement" \
  | table _time, host, Account_Name, src_ip, Logon_Type, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [9] T1087 - Account Discovery (Discovery)
# Common enumeration commands: whoami, net user, net group, quser.
# ---------------------------------------------------------------------
[SOCLab - T1087 - Account Discovery]
description = Reconnaissance commands to enumerate accounts (T1087).
enableSched = 1
cron_schedule = */5 * * * *
dispatch.earliest_time = -6m@m
dispatch.latest_time = -1m@m
search = index=sysmon EventCode=1 \
  (Image="*\\whoami.exe" OR Image="*\\net.exe" OR Image="*\\net1.exe" OR Image="*\\quser.exe") \
  CommandLine IN ("*whoami*", "*net user*", "*net group*", "*net localgroup*", "*quser*") \
  | eval mitre_technique="T1087", mitre_tactic="Discovery" \
  | table _time, host, User, Image, CommandLine, ParentImage, mitre_technique, mitre_tactic \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# ---------------------------------------------------------------------
# [10] T1566.001 + T1071.001 - Suspicious Outbound C2 Beacon
# Multi-technique: initial access follow-through + C2 web protocol.
# Uses Zeek conn.log for beaconing pattern (regular interval, small payload).
# ---------------------------------------------------------------------
[SOCLab - T1071.001 - HTTP Beaconing]
description = Beaconing pattern on outbound HTTP/HTTPS — regular intervals, low payload (T1071.001 + T1566.001).
enableSched = 1
cron_schedule = */10 * * * *
dispatch.earliest_time = -15m@m
dispatch.latest_time = -1m@m
search = index=zeek source="*conn.log" \
  | eval orig_bytes=tonumber(orig_bytes) \
  | bin _time span=30s \
  | stats count, avg(orig_bytes) as avg_bytes, dc(_time) as intervals by src_ip, dst_ip, dst_port \
  | where dst_port IN ("80", "443", "8080", "8443") AND intervals >= 5 AND (isnull(avg_bytes) OR avg_bytes < 2000) \
  | eval mitre_technique="T1071.001+T1566.001", mitre_tactic="Command and Control" \
  | collect index=notable
action.summary_index = 1
action.summary_index._name = notable

# =====================================================================
# Additional MITRE technique coverage (embedded in above searches):
# - T1055 (Process Injection): covered via Sysmon 8/10 in future rule set
# - T1098 (Account Manipulation): auditd useradd_exec/usermod_exec keys
# =====================================================================
